# alpine-consul-nodejs/root/etc/services.d/app/run
#!/usr/bin/env sh

# cd into our directory
cd /app

# start our node.js application
exec node server.js;
